# deploy/systemd/open-lzt-autoupdate.timer
[Unit]
Description=Run the open-lzt auto-updater periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
